# package/zl100mt/zl100mt-app/files/temperature.cron
*/1	*	*	*	*	root	/usr/bin/temperature.sh
*/1	*	*	*	*	root	sleep 30; /usr/bin/temperature.sh
