# package/zl100mt/zl100mt-app/files/etc/zl100mt-app/zl100mt-app.init
#!/bin/sh /etc/rc.common

START=80
STOP=90

start() {
    /usr/sbin/zl100mt-app &
}
stop() {
    killall -w zl100mt-app
}

