# target/linux/mvebu/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2014-2015 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

. /lib/functions/uci-defaults.sh
. /lib/mvebu.sh

ucidef_set_interface_loopback

board=$(mvebu_board_name)

case "$board" in
armada-385-linksys-caiman)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	;;
armada-385-linksys-cobra)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	;;
armada-xp-linksys-mamba)
	ucidef_set_interfaces_lan_wan "eth0" "eth1"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 5"
	ucidef_add_switch_vlan "switch0" "2" "4 6"
	;;
armada-385-db-ap)
	ucidef_set_interfaces_lan_wan "eth0 eth1" "eth2"
	;;
armada-xp-gp)
	ucidef_set_interface_lan "eth0 eth1 eth2 eth3"
	;;
rtrom01)
        #ucidef_set_interfaces_lan_wan "eth0 eth1" "eth2 usb0"
        ucidef_set_interface_lan "eth0"
        ucidef_set_interface_raw "lte_4g" "usb0" "ncm"
        uci batch <<EOF
set network.lte_4g='interface'
set network.lte_4g.ifname='usb0'
set network.lte_4g.proto='ncm'
set network.lte_4g.device='/dev/ttyUSB4'
set network.lte_4g.apn='unknown'
EOF
        ucidef_set_interface_raw "lte_z" "eth2" "dhcp"
        #ucidef_add_switch "switch0" "1" "0"
        #ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 5"
        #ucidef_add_switch_vlan "switch0" "2" "4 6"
        ;;
*)
	ucidef_set_interface_lan "eth0"
	;;
esac

uci commit network

exit 0
